# scripts/com.menubank.backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.menubank.dailybackup</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/musamecit/Desktop/QRMenu/scripts/daily-backup.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>0</integer>
        <key>Minute</key>
        <integer>1</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/tmp/menubank-backup.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/menubank-backup.err</string>
    <key>WorkingDirectory</key>
    <string>/Users/musamecit/Desktop/QRMenu</string>
</dict>
</plist>
